# crontab
5 4,9,14 * * * claude -p "Hi" >> /proc/1/fd/1 2>&1
